# debian/postinst
#!/bin/bash

echo "Reloading sysctl settings for inotify settings."
invoke-rc.d procps start
