# debian/postinst
#!/bin/sh

echo "Reloading sysctl settings for inotify settings."
invoke-rc.d procps start || true
